# central-hub/central-hub.service
[Unit]
Description=Central Hub Service - 网络基础设施中枢服务
Documentation=file:///vol1/1000/code/auto-ddnns/central-hub/README.md
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /vol1/1000/code/auto-ddnns/central-hub/server.mjs
WorkingDirectory=/vol1/1000/code/auto-ddnns/central-hub
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=CONFIG_PATH=/vol1/1000/code/auto-ddnns/central-hub/config/hub.json

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=central-hub

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/vol1/1000/code/auto-ddnns/central-hub/data /vol1/1000/code/auto-ddnns/central-hub/logs /vol1/1000/code/auto-ddnns/data

[Install]
WantedBy=default.target
