# scripts/sync-lucky-to-sunpanel.service
[Unit]
Description=Lucky to SunPanel Auto Sync
Documentation=file:///vol1/1000/code/auto-ddnns/docs/lucky-to-sunpanel.md
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/node /vol1/1000/code/auto-ddnns/scripts/sync-lucky-to-sunpanel.mjs --sync
WorkingDirectory=/vol1/1000/code/auto-ddnns
User=leecaiy
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
